# tmp/mcp-full.service
[Unit]
Description=Full MCP Scaffold Server
After=network.target

[Service]
Type=simple
User=mcp
WorkingDirectory=/opt/mcp/full
ExecStart=/usr/bin/node /opt/mcp/full/full-server.js
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
