# services/metrics-proxy/deploy/user.metrics-proxy.service
[Unit]
Description=Metrics Proxy (user)
After=default.target

[Service]
Type=simple
Environment=NODE_ENV=production
WorkingDirectory=%h/apps/metrics-proxy
ExecStart=%h/apps/metrics-proxy/start.sh
Restart=always
RestartSec=3

[Install]
WantedBy=default.target
