# systemd/mcp-tunnel-check.service
[Unit]
Description=MCP tunnel monitor (runs scripts/mcp_tunnel_check.sh) for MCP service
After=network.target

[Service]
Type=simple
ExecStart=/bin/bash /home/simon/Learning-Management-System-Academy/scripts/mcp_tunnel_check.sh
Restart=always
RestartSec=10
Environment=HOME=/home/simon
WorkingDirectory=/home/simon/Learning-Management-System-Academy

[Install]
WantedBy=default.target
