# .continue/systemd/agent-data_science.service
[Unit]
Description=Agent: data_science (runs the Data Science FastAPI agent)
After=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/simon/Learning-Management-System-Academy
ExecStartPre=/usr/bin/env bash -c 'mkdir -p /home/simon/.local/state/agents'
StandardOutput=append:/home/simon/.local/state/agents/data_science.log
StandardError=inherit
ExecStart=/usr/bin/env python3 /home/simon/Learning-Management-System-Academy/.continue/agents/agents_continue/data_science_server.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
