# .continue/systemd/agents.target
[Unit]
Description=Agents target grouping all example agent services
Wants=agent-core_dev.service agent-data_science.service agent-geo_intel.service agent-web_lms.service agent-systemops.service agent-legal_advisor.service agent-portfolio.service
After=network-online.target

[Install]
WantedBy=default.target
