# .continue/systemd/agent-web_lms.service
[Unit]
Description=Agent: web_lms (runs the Web & LMS FastAPI agent)
After=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/simon/Learning-Management-System-Academy
ExecStartPre=/usr/bin/env bash -c 'mkdir -p /home/simon/.local/state/agents'
StandardOutput=append:/home/simon/.local/state/agents/web_lms.log
StandardError=inherit
ExecStart=/usr/bin/env python3 /home/simon/Learning-Management-System-Academy/.continue/agents/agents_continue/web_lms_server.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
