# .continue/systemd/example-agent.service
[Unit]
Description=Example agent (connects to local SSE endpoint and logs model-list)
After=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/env python3 /home/simon/Learning-Management-System-Academy/.continue/agents/example_agent.py --once
Restart=on-failure
RestartSec=10

[Install]
WantedBy=default.target
