# .continue/systemd/health-check.timer
[Unit]
Description=Run agents health-check every 30 minutes (user)

[Timer]
OnBootSec=5min
OnUnitActiveSec=30min
Persistent=true

[Install]
WantedBy=timers.target
