# .continue/systemd/health-check.service
[Unit]
Description=Run agents health-check script (user)

[Service]
Type=oneshot
ExecStart=/bin/bash -lc '/home/simon/Learning-Management-System-Academy/.continue/scripts/health_check_agents.sh'
WorkingDirectory=/home/simon/Learning-Management-System-Academy
TimeoutStartSec=60
