# .continue/systemd/mcp-health.service
[Unit]
Description=MCP health check (runs .continue/scripts/mcp_health_check.sh)
Wants=mcp-agent.service mcp-tunnel.service

[Service]
Type=oneshot
ExecStart=/home/simon/Learning-Management-System-Academy/.continue/scripts/mcp_health_check.sh

[Install]
WantedBy=timers.target
