# .continue/systemd/mcp-tunnel.service
[Unit]
Description=SSH tunnel for Ollama MCP (user service)
After=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/ssh -i /home/simon/.ssh/id_ed25519_mcp -o ExitOnForwardFailure=yes -o ServerAliveInterval=60 -o ServerAliveCountMax=3 -N -L 11434:10.0.0.110:11434 -J root@136.243.155.166:2222 simonadmin@10.0.0.110
Restart=always
RestartSec=5
KillMode=process
StandardOutput=append:/home/simon/.local/state/mcp-tunnel.log
StandardError=append:/home/simon/.local/state/mcp-tunnel.err

[Install]
WantedBy=default.target
