# .continue/systemd/poll-to-sse.gunicorn.service
[Unit]
Description=Poll-to-SSE proxy (gunicorn)
After=network.target

# This unit is intended for system-level deployment; removed explicit User= to
# avoid invalid settings when installed as a user unit. If you want to run with
# a specific system user, install this under /etc/systemd/system/ and set User=.
[Service]
Type=simple
WorkingDirectory=/home/simon/Learning-Management-System-Academy/.continue/agents
# Ensure the environment is unbuffered for logs
Environment=PYTHONUNBUFFERED=1
# Use gunicorn to serve the Flask app defined in poll_to_sse.py (app variable)
ExecStart=/usr/bin/env gunicorn --workers 2 --threads 4 --bind 127.0.0.1:11534 poll_to_sse:app
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
