# .continue/systemd/smart-agent.service
[Unit]
Description=Smart agent (user) - run once to check system and report

[Service]
Type=oneshot
ExecStart=/usr/bin/python3 /home/simon/Learning-Management-System-Academy/.continue/agents/smart_agent.py --once --fix
WorkingDirectory=/home/simon/Learning-Management-System-Academy
Nice=10
TimeoutStartSec=120

[Install]
WantedBy=default.target
